# system_files/etc/systemd/system/framework-charge-limit.service
[Unit]
Description=Set Framework battery charge limit to 80%
After=multi-user.target
# Skip on hardware without this knob (VMs, non-Framework, BAT-index changes).
ConditionPathExists=/sys/class/power_supply/BAT1/charge_control_end_threshold

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/sh -c 'echo 80 > /sys/class/power_supply/BAT1/charge_control_end_threshold'
# Tolerate battery-driver init race after multi-user.target.
Restart=on-failure
RestartSec=5s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
