# system_files/etc/systemd/system/framework-charge-limit.service
[Unit]
Description=Set Framework battery charge limit to 80%
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=true
# Iterate all batteries; skip ones without the knob (VMs, non-Framework).
ExecStart=/bin/sh -c 'set -e; matched=0; for f in /sys/class/power_supply/BAT*/charge_control_end_threshold; do [ -w "$f" ] || continue; echo 80 > "$f"; matched=1; done; [ "$matched" = 1 ] || exit 0'
# Tolerate battery-driver init race after multi-user.target.
Restart=on-failure
RestartSec=5s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
