# system_files/etc/systemd/system/inhibit-sleep-when-docked.service
[Unit]
Description=Inhibit suspend while CalDigit TS4 is docked
# Lifecycle driven by 70-caldigit-ts4.rules (udev add/remove).

[Service]
Type=simple
# Cover lid-close too: TS4 doesn't always assert the ACPI dock bit that
# logind's HandleLidSwitchDocked=ignore relies on.
ExecStart=/usr/bin/systemd-inhibit --what=sleep:handle-lid-switch --who=blueberry --why="Docked to TS4" --mode=block /usr/libexec/blueberry/dock-monitor.sh
Restart=on-failure
