# system_files/etc/systemd/system/lock-before-sleep.service
[Unit]
Description=Lock sway sessions before suspend
Before=sleep.target

[Service]
Type=forking
ExecStart=/usr/bin/sh -c 'for u in $(loginctl list-users --no-legend | awk "{print \$2}"); do runuser -l "$u" -c "swaylock -f" || true; done'

[Install]
WantedBy=sleep.target
