# system_files/etc/systemd/system/usbguard.service.d/10-ensure-rules-file.conf
[Unit]
# rules.conf is populated post-install by `ujust setup-usbguard`. Daemon
# refuses to start without the file (and refuses world-readable rules), so
# pre-create an empty 0600 placeholder.

[Service]
ExecStartPre=/bin/sh -c 'test -f /etc/usbguard/rules.conf || install -m 0600 -o root -g root /dev/null /etc/usbguard/rules.conf'
