# system_files/etc/systemd/system/usbguard.service.d/10-ensure-rules-file.conf
[Unit]
# rules.conf is intentionally not shipped in the image — it is populated
# post-install via `ujust setup-usbguard` (which runs `usbguard generate-policy`).
# Without the file present, `usbguard-daemon` exits non-zero and systemd
# enters a restart loop. Create an empty placeholder with the correct mode
# (usbguard refuses to load rules from a world-readable file) so the daemon
# can start; rules will be filled in by the setup recipe.

[Service]
ExecStartPre=/bin/sh -c 'test -f /etc/usbguard/rules.conf || install -m 0600 -o root -g root /dev/null /etc/usbguard/rules.conf'
